# root/etc/services.d/st2/run
#!/usr/bin/with-contenv sh

umask 022

echo "Executing => ST2"
exec \
s6-setuidgid abc /app/supertransfer2.sh